# files/etc/systemd/user/podman-auto-update.timer.d/override.conf
[Timer]
RandomizedDelaySec=15m
OnBootSec=15m
OnCalendar=*-*-* 04:00:00
Persistent=true
OnUnitInactiveSec=12h
