# files/etc/profile.d/glorpfetch.sh
alias glorpfetch="/usr/bin/fastfetch --logo /etc/fastfetch/glorp.txt"
